# raspberry/pi-online-notify.service
[Unit]
Description=Raspberry Pi online notifier (Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=faust
EnvironmentFile=-/etc/ege-app/pi-notify.env
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash /usr/local/bin/pi-online-notify.sh
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
