# raspberry/pocketbase.service
[Unit]
Description=PocketBase Service for EGE Tasks Manager
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=faust
Group=faust
WorkingDirectory=/home/faust/ege-app/pocketbase
ExecStart=/home/faust/ege-app/pocketbase/pocketbase serve --http=127.0.0.1:8090
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов (для Raspberry Pi)
MemoryLimit=512M
CPUQuota=80%

# Переменные окружения
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
